# deploy/systemd/peersunity-reverb.service
[Unit]
Description=PeersUnity Laravel Reverb WebSocket Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/peersunity.com
ExecStart=/usr/bin/php /var/www/peersunity.com/artisan reverb:start --host=127.0.0.1 --port=8080 --no-interaction
ExecStop=/usr/bin/php /var/www/peersunity.com/artisan reverb:restart
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=append:/var/www/peersunity.com/storage/logs/reverb.log
StandardError=append:/var/www/peersunity.com/storage/logs/reverb.log
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
